# radres/deploy/supervisor.ini
[program:radsfile]
command=/home/swihart/supervisor/scripts/start_radsfile
directory=/home/swihart/webapps/radsfile
user=swihart
autostart=true
autorestart=true
redirect_stderr=True
